# supervisor/database.conf
[supervisord]
nodaemon = true

[program:database]
; Database Daemon (Port 27017)
command = /usr/bin/mongod --config /etc/mongodb.conf
user = root
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/database.log
stderr_logfile = /var/log/supervisor/database.log
